# product/package/batocera-luigios-dev-connectivity/files/package-files.manifest
/etc/init.d/S07luigios-dev-connectivity
/usr/libexec/luigios/dev-connectivity-guard
/usr/share/luigios/manifests/dev-connectivity.files
